# mcp/pacemaker.service.in
[Unit]
Description=Pacemaker High Availability Cluster Manager
After=network.target
Requires=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
KillMode=process
NotifyAccess=main
SysVStartPriority=99
EnvironmentFile=-@sysconfdir@/sysconfig/pacemaker

ExecStart=@sbindir@/pacemakerd -f

# Pacemaker can only exit after all managed services have shut down
# A HA database could conceivably take even longer than this 
TimeoutSec=30min

# Restart options include: no, on-success, on-failure, on-abort or always
Restart=on-failure
