# cts/pacemaker-cts-dummyd@.service.in
[Unit]
Description=Dummy daemon for Pacemaker CTS testing

[Service]
Type=notify
ExecStart=@CRM_DAEMON_DIR@/pacemaker-cts-dummyd %i

[Install]
DefaultInstance=0
