# daemons/execd/pacemaker_remote.service.in
[Unit]
Description=Pacemaker Remote executor daemon
Documentation=man:pacemaker-remoted http://clusterlabs.org/doc/en-US/Pacemaker/1.1-pcs/html/Pacemaker_Remote/index.html

# See main pacemaker unit file for descriptions of why these are needed
After=network.target
After=time-sync.target
After=dbus.service
Wants=dbus.service
After=resource-agents-deps.target
Wants=resource-agents-deps.target
After=syslog.service
After=rsyslog.service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
KillMode=process
NotifyAccess=none
EnvironmentFile=-@CONFIGDIR@/pacemaker
EnvironmentFile=-@CONFIGDIR@/sbd

ExecStart=@sbindir@/pacemaker-remoted

# Systemd v227 and above can limit the number of processes spawned by a
# service. That is a bad idea for an HA cluster resource manager, so disable it
# by default. The administrator can create a local override if they really want
# a limit. If your systemd version does not support TasksMax, and you want to
# get rid of the resulting log warnings, comment out this option.
TasksMax=infinity

# Pacemaker Remote can exit only after all managed services have shut down;
# an HA database could conceivably take even longer than this 
TimeoutStopSec=30min
TimeoutStartSec=30s

# Restart options include: no, on-success, on-failure, on-abort or always
Restart=on-failure

# crm_perror() writes directly to stderr, so ignore it here
# to avoid double-logging with the wrong format
StandardError=null
