# tools/crm_mon.service.in
[Unit]
Description=Daemon for pacemaker monitor

[Service]
Type=forking
EnvironmentFile=-@sysconfdir@/sysconfig/crm_mon
ExecStart=@sbindir@/crm_mon $OPTIONS
Restart=always

# crm_perror() writes directly to stderr, so ignore it here to avoid
# double-logging with the wrong format
StandardError=null

[Install]
WantedBy=multi-user.target
